# systemd/weight-monitor.service.template
[Unit]
Description=Pet feeder weight jam monitor
After=network-online.target time-sync.target
Wants=network-online.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKDIR%
ExecStart=%WORKDIR%/.venv/bin/wm-daemon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
